# scripts/systemd/far-ollama-selective-acceptance.service
[Unit]
Description=FAR D-backed Ollama for P14 selective acceptance

[Service]
Type=simple
WorkingDirectory=/mnt/d/FAR-workspace/FAR-longterm
ExecStart=/bin/bash -lc 'source ~/miniconda3/etc/profile.d/conda.sh; conda activate train; source scripts/windows_gpu_env.sh; exec ollama serve'
Restart=always
RestartSec=5s
StandardOutput=append:/mnt/d/FAR-outputs/selective_acceptance_v2.ollama.log
StandardError=append:/mnt/d/FAR-outputs/selective_acceptance_v2.ollama.log

[Install]
WantedBy=default.target
