# scripts/systemd/far-ollama-p5.service
[Unit]
Description=FAR D-backed Ollama service for registered P5 ablations

[Service]
Type=simple
WorkingDirectory=/mnt/d/FAR-workspace/FAR-longterm
ExecStart=/bin/bash -lc 'source ~/miniconda3/etc/profile.d/conda.sh; conda activate train; source scripts/windows_gpu_env.sh; exec ollama serve'
Restart=always
RestartSec=5s
StandardOutput=append:/mnt/d/FAR-outputs/p5_ramdocs_v1.ollama.log
StandardError=append:/mnt/d/FAR-outputs/p5_ramdocs_v1.ollama.log

[Install]
WantedBy=default.target
